# k8s-secret-chowner.conf
[program:k8s-secret-chowner]
command=/usr/local/libexec/k8s-secret-chowner.sh
startsecs=0

# Log to stdout
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
